# supervisor/plotter.conf
[program:api.plotter]
command = python plotter_srv.py
directory=/opt/plotter/bin/
environment=PATH="/opt/plotter/bin/venv/bin:%(ENV_PATH)s"

numprocs=1
numprocs_start=0

logfile_maxbytes = 20MB
logfile_backups = 10

autorestart=true
startsecs=0

user=root